# avatar-lipsync.service
[Unit]
Description=Avatar Lipsync Web Service
After=network.target

[Service]
Type=simple
User=arman
WorkingDirectory=/home/arman/workspace2
Environment="PATH=/home/arman/workspace2/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/arman/workspace2/.venv/bin/python /home/arman/workspace2/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Ограничения ресурсов (опционально)
# MemoryLimit=8G
# CPUQuota=400%

[Install]
WantedBy=multi-user.target
