# avatar-lipsync.service
[Unit]
Description=Avatar Lipsync Web Service
After=network.target

[Service]
Type=simple
User=arman
WorkingDirectory=/home/arman/musetalk/avatar
Environment="PATH=/home/arman/musetalk/avatar/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/arman/musetalk/avatar/venv/bin/python /home/arman/musetalk/avatar/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Ограничения ресурсов (опционально)
# MemoryLimit=8G
# CPUQuota=400%

[Install]
WantedBy=multi-user.target
